# service/com.victronenergy.mppsolar.service
[Unit]
Description=MPP Solar D-Bus service for Venus OS
After=dbus.service
Requires=dbus.service

[Service]
Type=simple
ExecStart=/opt/victronenergy/dbus-mppsolar/dbus-mppsolar.py
Restart=always
RestartSec=10
User=root
Group=root
Environment=PYTHONPATH=/opt/victronenergy/dbus-mppsolar

[Install]
WantedBy=multi-user.target